# elk/curator/crontab
@daily root /usr/local/bin/curator /root/.curator/curator_action.yml
# An empty line is required at the end of this file for a valid cron file.
